# scripts/scaleway-deploy/dockly-auto-shutdown.timer
[Unit]
Description=Check DocklyOCR GPU idle every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=30s

[Install]
WantedBy=timers.target
