# scripts/scaleway-deploy/dockly-auto-shutdown.timer
[Unit]
Description=Check DocklyOCR GPU idle every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
